# redhat/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local
# Start the VAP Demo Node Server
bash /home/ec2-user/vap/redhat/start_vap.sh > /var/log/vap.rc.local.log 2>&1
